# deploy/feishu-claude.service
[Unit]
Description=Feishu Claude Code Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/feishu-claude
ExecStart=/path/to/feishu-claude/.venv/bin/python main.py
Restart=always
RestartSec=10

# 环境变量（也可以用 EnvironmentFile 指向 .env）
EnvironmentFile=/path/to/feishu-claude/.env

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
